# scripts/startup.service
[Unit]
Description=Startup Service

[Service]
ExecStart=/home/tony/.scripts/startup.sh

[Install]
WantedBy=default.target
