# scripts/startup.service
[Service]
Type=oneshot
RemainAfterExit=true
StandardOutput=journal
ExecStart=/home/tony/.scripts/startup.sh
[Install]
WantedBy=default.target
